# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── AADS API 서버 ──────────────────────────────────────────────────
[program:aads-api]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/aads-api.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/aads-api.log
stderr_logfile_maxbytes=10MB
priority=10

; ── MCP Filesystem 서버 ─────────────────────────────────────────────
[program:mcp-filesystem]
command=python -m mcp_servers.filesystem_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-filesystem.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-filesystem.log
stderr_logfile_maxbytes=5MB
priority=5
environment=MCP_FS_ROOT="/tmp/aads_workspace"

; ── MCP Git 서버 ────────────────────────────────────────────────────
[program:mcp-git]
command=python -m mcp_servers.git_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-git.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-git.log
stderr_logfile_maxbytes=5MB
priority=5
environment=MCP_GIT_ROOT="/tmp/aads_workspace"

; ── MCP Memory 서버 ─────────────────────────────────────────────────
[program:mcp-memory]
command=python -m mcp_servers.memory_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-memory.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-memory.log
stderr_logfile_maxbytes=5MB
priority=5

[group:mcp-servers]
programs=mcp-filesystem,mcp-git,mcp-memory
