# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; в”Җв”Җ AADS API м„ңлІ„ в”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җ
[program:aads-api]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1 --timeout-keep-alive 120 --timeout-graceful-shutdown 120
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=SIGTERM
stopwaitsecs=120
stdout_logfile=/var/log/aads-api.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/aads-api.log
stderr_logfile_maxbytes=10MB
priority=10

; в”Җв”Җ MCP Filesystem м„ңлІ„ в”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җ
[program:mcp-filesystem]
command=python -m mcp_servers.filesystem_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-filesystem.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-filesystem.log
stderr_logfile_maxbytes=5MB
priority=5
environment=MCP_FS_ROOT="/tmp/aads_workspace"

; в”Җв”Җ MCP Git м„ңлІ„ в”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җ
[program:mcp-git]
command=python -m mcp_servers.git_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-git.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-git.log
stderr_logfile_maxbytes=5MB
priority=5
environment=MCP_GIT_ROOT="/root/aads/aads-server"

; в”Җв”Җ MCP Memory м„ңлІ„ в”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җ
[program:mcp-memory]
command=python -m mcp_servers.memory_server
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/mcp-memory.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp-memory.log
stderr_logfile_maxbytes=5MB
priority=5

; в”Җв”Җ Playwright MCP м„ңлІ„ (AADS-159) в”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җв”Җ
; Python playwrightк°Җ in-process(aads-api)лЎң мӢӨн–үлҗҳлҜҖлЎң autostart=false
; Node.js @playwright/mcpлҘј лі„лҸ„ мӮ¬мҡ© мӢң nodejs м„Өм№ҳ + autostart=true ліҖкІҪ
[program:playwright-mcp]
command=npx @playwright/mcp@latest --headless --port 8768
directory=/app
autostart=false
autorestart=false
startsecs=5
startretries=1
stdout_logfile=/var/log/playwright-mcp.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/playwright-mcp.log
stderr_logfile_maxbytes=5MB
priority=5
environment=PLAYWRIGHT_BROWSERS_PATH="/root/.cache/ms-playwright"

[group:mcp-servers]
programs=mcp-filesystem,mcp-git,mcp-memory,playwright-mcp
