# scripts/aads-pipeline-runner.service
[Unit]
Description=AADS Pipeline Runner — DB-driven Claude Code executor
After=docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/aads/aads-server
ExecStart=/root/aads/aads-server/scripts/pipeline-runner.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aads-pipeline-runner

# 민감정보는 별도 파일에서 로드 (H2: 하드코딩 제거)
EnvironmentFile=-/root/.config/aads-runner.env
Environment=POLL_INTERVAL=5
Environment=MAX_RUNTIME=7200
Environment=MAX_RETRIES=2
Environment=RUNNER_PROJECTS=AADS
Environment=HOME=/root
Environment=PATH=/root/.nvm/versions/node/v20.20.0/bin:/usr/local/bin:/usr/bin:/bin

LimitNOFILE=65535
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
